# contrib/init/ravend.conf
description "Raven Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env AVIAND_BIN="/usr/bin/aviand"
env AVIAND_USER="raven"
env AVIAND_GROUP="raven"
env AVIAND_PIDDIR="/var/run/aviand"
# upstart can't handle variables constructed with other variables
env AVIAND_PIDFILE="/var/run/aviand/aviand.pid"
env AVIAND_CONFIGFILE="/etc/raven/raven.conf"
env AVIAND_DATADIR="/var/lib/aviand"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # aviand will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$AVIAND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run aviand."
        echo "The setting must appear in $AVIAND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Raven Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$AVIAND_PIDDIR"
    chmod 0755 "$AVIAND_PIDDIR"
    chown $AVIAND_USER:$AVIAND_GROUP "$AVIAND_PIDDIR"
    chown $AVIAND_USER:$AVIAND_GROUP "$AVIAND_CONFIGFILE"
    chmod 0660 "$AVIAND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$AVIAND_PIDFILE" \
    --chuid $AVIAND_USER:$AVIAND_GROUP \
    --exec "$AVIAND_BIN" \
    -- \
    -pid="$AVIAND_PIDFILE" \
    -conf="$AVIAND_CONFIGFILE" \
    -datadir="$AVIAND_DATADIR" \
    -disablewallet \
    -daemon

